# service/livesport-crawler.service
[Unit]
Description=Script for crawling of livescore.com site

[Service]
Type=simple
User=marxin
Group=marxin
ExecStart=/home/marxin/Programming/livesport-crawler/target/debug/livesport-crawler https://www.livesport.cz/tym/sparta-praha/zcG9U7N6/ "Sparta Praha" /var/www/html/sparta.json firefox -k
Restart=on-failure
Wants=network-online.target
After=network-online.target
RestartSec=60

[Install]
WantedBy=default.target
